# src/nix/root/etc/s6-overlay/s6-rc.d/init-feature-nix-add-packages/up
/etc/s6-overlay/s6-rc.d/init-feature-nix-add-packages/run